# profiles/dev/rc.local
#!/bin/bash
usermod -p '' devuser
passwd --expire devuser
rm -f /etc/rc.local
exit 0
