# profiles/sec/rc.local
#!/bin/bash
usermod -p '' secretary
passwd --expire secretary
rm -f /etc/rc.local
exit 0
