# profiles/admin/rc.local
#!/bin/bash
usermod -p '' adminsys
passwd --expire adminsys
rm -f /etc/rc.local
exit 0
