# profiles/admin/rc.local
#!/bin/bash
usermod -p '' admin
passwd --expire admin
rm -f /etc/rc.local
exit 0
